# metrics/api/collectd/remove-collectd-csv.cron
SHELL=/bin/bash
PATH=/usr/bin:/bin
0 0 * * * find /var/lib/collectd/csv/ -type f -mtime +1 -exec rm -vf {} \; 2>&1 | logger -t 'Metrics API'
